# DEBIAN/prerm
#!/bin/bash
set -e

SERVICE_NAME="user-activity-monitor.service"

# Stop and disable the service
systemctl stop $SERVICE_NAME || true
systemctl disable $SERVICE_NAME || true

# Remove the service file
rm -f /etc/systemd/system/$SERVICE_NAME

# Reload systemd
systemctl daemon-reload
